# MCT_Referee/init_files/mct-registry.conf
description "The MultiClouds Tournament Registry Service"
author "mrmassis"

start on runlevel [2345]
stop on runlevel [!2345]


respawn

exec /usr/lib/python2.7/dist-packages/mct/MCT_Registry.py
